# luci/themes/luci-theme-darky/root/etc/uci-defaults/30_luci-theme-darky
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.Darky=/luci-static/darky
		set luci.main.mediaurlbase=/luci-static/darky
		commit luci
	EOF
fi

exit 0
